# files/deploy.cron
5 * * * * root /var/www/html/TWLight/bin/twlight_deploy.sh >/dev/null 2>&1
